# systemd/astrowidget-grade.service
[Unit]
Description=astrowidget FITS auto-grader (Phase 3 transparency ground truth)
Documentation=https://github.com/dustinspace217/astrowidget
Wants=network-online.target
After=network-online.target
# The sweep reads FITS off the NAS, so the CIFS mount must be up first. Naming the
# full data path lets systemd resolve whichever mount actually covers it (the mount
# point may be /mnt/vaultinspace or something deeper) and wait for it. If the NAS is
# still unreachable, the grader fails loud: main() exits non-zero when the root isn't
# a readable directory, so the unit shows failed rather than silently grading nothing.
RequiresMountsFor=/mnt/vaultinspace/Astrophotography/Raws
# On failure (non-zero exit: NAS unreadable, unhandled error), fire a desktop
# notification so an unattended failure is visible, not just red in systemctl.
OnFailure=astrowidget-grade-failure.service

[Service]
# Runs once per timer fire, then exits (like the fetcher). It walks last night's
# session folders, grades the not-yet-graded ones, writes fits_grades, and stops.
Type=oneshot
# --scan: sweep the Raws/ tree and grade every complete, ungraded night within the
# lookback window. --rig "Eon 70" restricts grading to the home scope — Raws/ also
# holds iTelescope rentals (T24/T26/…), which must NOT be logged as Bainbridge.
# Coordinates for dawn-exclusion are read from the gitignored config.toml at runtime
# (so no lat/lon is baked into this unit). --since-days defaults to 30; the first run
# after a long gap can backfill more with a manual
#   python3 grader/grade.py --scan <root> --site Bainbridge --rig "Eon 70" --since-days 0 --write
ExecStart=/usr/bin/python3 %h/Claude/astrowidget/grader/grade.py --scan /mnt/vaultinspace/Astrophotography/Raws --site Bainbridge --rig "Eon 70" --write
StandardOutput=journal
StandardError=journal
# Reads one sub into memory at a time (read_sub returns scalar metrics, not the pixel
# array). Steady-state nightly runs are light, but the one-time `--since-days 0`
# backfill over many sessions peaked at ~1.2 GB and spilled to swap, so the cap is 2 GB
# for headroom (still catches a runaway without disturbing the desktop).
MemoryMax=2G
# A first-run backfill over CIFS can read a lot of subs. Cap at 30 min so a hung
# network read can't wedge the timer indefinitely; if a genuine backfill needs
# longer, run it manually once (it's idempotent — reruns skip graded nights).
TimeoutStartSec=1800
